# pic.service
[Unit]
Description = pic daemon

[Service]
ExecStart = /opt/python_pic/run.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
